# systemd/backup-minecraft.timer
[Unit]
Description=Daily Minecraft Server backup

[Timer]
OnCalenday=daily

[Install]
WantedBy=timers.target
