# systemd/backup-minecraft.timer
[Unit]
Description=Daily Minecraft Server backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
